# yi-hack-v3+/yi-hack-v3/etc/crontabs/root
0 * * * * /home/yi-hack/script/clean_records.sh 2
